# deploy/alpaca-stream-monitor.service
[Unit]
Description=Alpaca WebSocket stream trading monitor
Wants=network-online.target
After=network-online.target
Conflicts=alpaca-bot.service

[Service]
Type=simple
WorkingDirectory=/opt/alpaca-bot
EnvironmentFile=/opt/alpaca-bot/.env
ExecStart=/opt/alpaca-bot/venv/bin/python -m src.stream_monitor
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
